# zuul/tomcat.conf
[program:tomcat]
command=/opt/tomcat/fixaddress_and_start.sh
