# cassandra/cassandra.conf
[program:cassandra]
command=/opt/cassandra/fixaddress_and_start.sh

